# templates/kibana.service.j2
[Unit]
Description=Elasticsearch
After=network.target

[Service]
Type=simple
ExecStart=/opt/kibana-4.5.1-linux-x64/bin/kibana
Restart=always
User=pi

[Install]
WantedBy=multi-user.target
